# connectedhomeip/connectedhomeip/third_party/silabs/gecko_sdk/util/third_party/ot-br-posix/src/web/otbr-web.service.in
[Unit]
Description=Border Router Web
After=otbr-agent.service
ConditionPathExists=@CMAKE_INSTALL_FULL_SBINDIR@/otbr-web

[Service]
EnvironmentFile=-@CMAKE_INSTALL_FULL_SYSCONFDIR@/default/otbr-web
ExecStart=@CMAKE_INSTALL_FULL_SBINDIR@/otbr-web $OTBR_WEB_OPTS
Restart=on-failure
RestartSec=5
RestartPreventExitStatus=SIGKILL

[Install]
WantedBy=multi-user.target
Alias=otbr-web.service
