# connectedhomeip/connectedhomeip/third_party/openthread/ot-efr32/third_party/silabs/gecko_sdk/util/third_party/ot-br-posix/src/agent/otbr-agent.service.in
[Unit]
Description=OpenThread Border Router Agent
ConditionPathExists=@CMAKE_INSTALL_FULL_SBINDIR@/otbr-agent
Requires=dbus.socket
After=dbus.socket

[Service]
EnvironmentFile=-@CMAKE_INSTALL_FULL_SYSCONFDIR@/default/otbr-agent
@EXEC_START_PRE@ExecStart=@CMAKE_INSTALL_FULL_SBINDIR@/otbr-agent $OTBR_AGENT_OPTS
KillMode=mixed
Restart=on-failure
RestartSec=5
RestartPreventExitStatus=SIGKILL

[Install]
WantedBy=multi-user.target
Alias=otbr-agent.service
